# rpi/powerblock-reboot-intent-dropin.service
[Service]
ExecStopPost=/usr/local/bin/powerblock-send-reboot-intent-if-reboot.sh
