# rpi/powerblock-reboot-intent.service
[Unit]
Description=PowerBlock reboot-intent hook
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/usr/local/bin/powerblock-send-reboot-intent.sh

[Install]
WantedBy=multi-user.target
